# mtk-openwrt-4.0.1.0/files/etc/init.d/startup-config
#!/bin/sh /etc/rc.common

# This script is startup-config using uci.
# The following commands overwrite default uci config.

START=99


stop() {
	echo "do noting..."
}

start() {
	uci set firewall.@zone[1].input='ACCEPT'
	uci set firewall.@zone[1].forward='ACCEPT'

	uci commit firewall
	/etc/init.d/firewall restart

	uci set network.lan.igmp_snooping='1'

	uci commit network
#	/etc/init.d/network restart
}
